# deploy/midia-indoor.service
# =============================================================
# Serviço systemd — Nexo Mídia (Gunicorn)
# Instalado automaticamente por deploy/scripts/install.sh em
# /etc/systemd/system/midia-indoor.service
#
# Comandos úteis:
#   sudo systemctl status midia-indoor
#   sudo systemctl restart midia-indoor
#   sudo journalctl -u midia-indoor -f
#
# Layout de pasta única: __APP_DIR__ é o próprio código (git ou
# cópia do zip), com venv/.env/uploads/logs/instance dentro dela.
# O install.sh substitui __APP_DIR__ pelo caminho real ao copiar
# este arquivo para /etc/systemd/system/.
# =============================================================
[Unit]
Description=Nexo Midia - Flask + Gunicorn (mídia indoor digital)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=midia-indoor
Group=midia-indoor
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/venv/bin/gunicorn -c deploy/gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=mixed

# ---- Hardening (defesa em profundidade a nível de SO) ----
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_DIR__

[Install]
WantedBy=multi-user.target
