# deploy/midia-indoor.service
# =============================================================
# Serviço systemd — Nexo Mídia (Gunicorn)
# Instalado automaticamente por deploy/scripts/install.sh em
# /etc/systemd/system/midia-indoor.service
#
# Comandos úteis:
#   sudo systemctl status midia-indoor
#   sudo systemctl restart midia-indoor
#   sudo journalctl -u midia-indoor -f
#
# Este arquivo usa o layout "releases + current" (estilo blue/green):
#   /opt/midia-indoor/current -> aponta para a release ativa
#   /opt/midia-indoor/shared  -> .env, venv, uploads, logs e banco
#     (compartilhados entre releases, preservados nas atualizações)
# =============================================================
[Unit]
Description=Nexo Midia - Flask + Gunicorn (mídia indoor digital)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=midia-indoor
Group=midia-indoor
WorkingDirectory=/opt/midia-indoor/current
EnvironmentFile=/opt/midia-indoor/shared/.env
ExecStart=/opt/midia-indoor/shared/venv/bin/gunicorn -c deploy/gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=mixed

# ---- Hardening (defesa em profundidade a nível de SO) ----
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/midia-indoor/shared

[Install]
WantedBy=multi-user.target
